# heatmiser-parse.service
[Service]
Type=oneshot
ExecStart=/opt/heatmiser/build_dataset.py