# heatmiser-parse.timer
[Timer]
OnCalendar=*:1,16,31,46,59

[Install]
WantedBy=multi-user.target