# heatmiser-parse.service
[Service]
Type=oneshot
ExecStart=/usr/bin/python /opt/heatmiser/build_dataset.py