# relay-poll.timer
[Timer]
OnCalendar=*-*-* *:*:22
AccuracySec=1ms

[Install]
WantedBy=multi-user.target