# relay-poll.service
[Service]
Type=oneshot
ExecStart=/opt/heatmiser/poll_relays.py