# relay-poll.service
[Service]
Type=oneshot
ExecStart=/usr/bin/python /opt/heatmiser/poll_relays.py